# templates/pgbouncer_exporter_healer.service.j2
[Unit]
Description=Pgbouncer Exporter Healer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pgbouncer_exporter_healer.sh